# agent/trpc_agent/template/supervisord.conf
[supervisord]
nodaemon=true

[program:server]
command=bun run src/index.js
directory=/app/server

[program:client]
command=caddy run --config /app/client/Caddyfile
directory=/app/client